# docker/spring/etc/supervisor/conf.d/envoy.conf
[program:envoy]
command=/usr/bin/envoy -c /etc/envoy.yaml -l debug
autostart=true
autorestart=true
user=root
