# docker/spring/etc/supervisor/conf.d/mysql.conf
[program:mysql]
command=/usr/bin/mariadbd-safe --datadir='/var/lib/mysql' --user=mysql
autostart=true
autorestart=true
user=root
