# docker/spring/etc/supervisor/conf.d/opensearch.conf
[program:opensearch]
command=/opt/opensearch-2.3.0/bin/opensearch
autostart=true
autorestart=true
user=nobody
directory=/opt/opensearch-2.3.0
environment=OPENSEARCH_JAVA_HOME="/opt/opensearch-2.3.0/jdk"
