# docker/spring/etc/supervisor/conf.d/opensearch.conf
[program:opensearch]
command=/opt/opensearch/bin/opensearch
autostart=true
autorestart=true
user=nobody
directory=/opt/opensearch
environment=OPENSEARCH_JAVA_HOME="/opt/opensearch/jdk"
