# docker/spring/etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=/usr/bin/php-fpm --nodaemonize --fpm-config /etc/php/php-fpm.conf
autostart=true
autorestart=true
user=root
