# docker/spring/etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm8.1]
command=/usr/sbin/php-fpm8.1 --nodaemonize --fpm-config /etc/php/8.1/fpm/php-fpm.conf
autostart=true
autorestart=true
user=root
