# docker/spring/etc/supervisor/conf.d/postgresql.conf
[program:postgresql]
command=/usr/bin/postgres -D /var/lib/postgres/data
autostart=true
autorestart=true
user=postgres
