# docker/spring/etc/supervisor/conf.d/postgresql.conf
[program:postgresql]
command=/usr/lib/postgresql/14/bin/postgres -D /var/lib/postgresql/data
autostart=true
autorestart=true
user=postgres
