# docker/jammy/etc/supervisor/conf.d/mysql.conf
[program:mysql]
command=/usr/bin/mysqld_safe --user=mysql
autostart=true
autorestart=true
user=root
