# docker/jammy/etc/supervisor/conf.d/rabbitmq.conf
[program:rabbitmq]
command=/usr/sbin/rabbitmq-server
autostart=true
autorestart=true
user=root
