# dashboard-v4/docker/spring/etc/supervisor/conf.d/redis.conf
[program:redis-s]
command=/usr/bin/redis-server /etc/redis/node-s.conf
autostart=true
autorestart=true
user=redis

[program:redis-1]
command=/usr/bin/redis-server /etc/redis/node-1.conf
autostart=true
autorestart=true
user=redis

[program:redis-2]
command=/usr/bin/redis-server /etc/redis/node-2.conf
autostart=true
autorestart=true
user=redis

[program:redis-3]
command=/usr/bin/redis-server /etc/redis/node-3.conf
autostart=true
autorestart=true
user=redis

[program:redis-4]
command=/usr/bin/redis-server /etc/redis/node-4.conf
autostart=true
autorestart=true
user=redis

[program:redis-5]
command=/usr/bin/redis-server /etc/redis/node-5.conf
autostart=true
autorestart=true
user=redis

[program:redis-6]
command=/usr/bin/redis-server /etc/redis/node-6.conf
autostart=true
autorestart=true
user=redis


