# deploy/roles/mint-v2/templates/v2/discussion-worker.service.j2
[Unit]
Description=Mint discussion worker for {{ app_deploy_env }}
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/php artisan mq:discussion
WorkingDirectory={{ app_deploy_root }}/htdocs
User=www-data
Group=www-data
Restart=always
RestartSec=30s

[Install]
WantedBy=multi-user.target
