# deploy/roles/mint-v2/templates/v2/morus.service.j2
[Unit]
Description=Mint morus rpc service for {{ app_deploy_env }}
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/php server.php -p {{ app_morus_listen_port }}
WorkingDirectory={{ app_deploy_root }}/agile/rpc/morus/morus
User=www-data
Group=www-data
Restart=always
RestartSec=30s

[Install]
WantedBy=multi-user.target
