# deploy/roles/mint-v2/templates/v2/scheduler.timer.j2
[Unit]
Description=Runs Mint mint-{{ app_deploy_env }} scheduler every minute

[Timer]
OnUnitActiveSec=1minute
Unit=mint-{{ app_deploy_env }}-scheduler.service

[Install]
WantedBy=multi-user.target
